# conf/octane/octane.ini
[program:octane]
command = /usr/local/bin/php /var/www/html/artisan octane:start --max-requests=%(ENV_OCTANE_MAX_REQUESTS)s --workers=%(ENV_OCTANE_WORKERS)s --task-workers=%(ENV_OCTANE_TASKS)s
numprocs = 1
autostart = %(ENV_OCTANE_START)s
autorestart = true
startsecs = 3
stdout_logfile = /var/www/html/storage/logs/supervisor/app_out.log
stderr_logfile = /var/www/html/storage/logs/supervisor/app_err.log